# config/systemd/chickencoop-monitor.service
[Unit]
Description=Chicken Coop Monitoring Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/chickencoop
Environment="COOP_ID=default"
ExecStartPre=/bin/bash -c 'sleep 5 && find /home/pi/chickencoop/videos/tmp -type f -mmin +60 -delete'
ExecStart=/home/pi/chickencoop/venv/bin/python -m src.main
Restart=always
RestartSec=10
StandardOutput=append:/home/pi/chickencoop/logs/coop-stdout.log
StandardError=append:/home/pi/chickencoop/logs/coop-stderr.log

[Install]
WantedBy=multi-user.target
